# packaging/arch/stophammer-node.install
post_install() {
  cat <<'EOF'
stophammer-node installed.

Next steps:
  1. Edit /etc/stophammer/community.env
  2. Enable the service:
       systemctl enable --now stophammer-community.service

This package intentionally does not include primary/indexer maintenance tools.
EOF
}

post_upgrade() {
  post_install
}
