# packaging/arch/stophammer-indexer.install
post_install() {
  cat <<'EOF'
stophammer-indexer installed.

Next steps:
  1. Edit /etc/stophammer/primary.env
  2. Enable the service:
       systemctl enable --now stophammer-primary.service

This package conflicts with stophammer-node by design.
Install stophammer-crawler separately if this host also runs feed ingestion.
EOF
}

post_upgrade() {
  post_install
}
