# packaging/systemd/stophammer-gossip.service
[Unit]
Description=Stophammer gossip crawler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stophammer-crawler
Group=stophammer-crawler
SupplementaryGroups=podping
EnvironmentFile=/etc/stophammer/crawler-gossip.env
ExecStart=/bin/sh -lc 'set -- /usr/bin/stophammer-crawler gossip --state "${GOSSIP_STATE_DB:-/var/lib/stophammer-crawler/gossip_state.db}" --skip-db "${FEED_SKIP_DB:-/var/lib/stophammer-crawler/feed_skip.db}"; if [ -n "${GOSSIP_SSE_URL:-}" ]; then set -- "$@" --sse-url "$GOSSIP_SSE_URL"; fi; if [ -n "${GOSSIP_ARCHIVE_DB:-}" ]; then set -- "$@" --archive-db "$GOSSIP_ARCHIVE_DB"; fi; if [ -n "${GOSSIP_SINCE_HOURS:-}" ]; then set -- "$@" --since-hours "$GOSSIP_SINCE_HOURS"; fi; if [ -n "${GOSSIP_SKIP_TTL_DAYS:-}" ]; then set -- "$@" --skip-ttl-days "$GOSSIP_SKIP_TTL_DAYS"; fi; if [ -n "${CONCURRENCY:-}" ]; then set -- "$@" --concurrency "$CONCURRENCY"; fi; if [ "${GOSSIP_SKIP_KNOWN_NON_MUSIC:-1}" = "1" ]; then set -- "$@" --skip-known-non-music; fi; if [ "${GOSSIP_QUIET:-0}" = "1" ]; then set -- "$@" --quiet; fi; exec "$@"'
WorkingDirectory=/var/lib/stophammer-crawler
Restart=on-failure
RestartSec=5
StateDirectory=stophammer-crawler
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
ReadWritePaths=/var/lib/stophammer-crawler

[Install]
WantedBy=multi-user.target
