# packaging/systemd/stophammer-import.service
[Unit]
Description=Stophammer one-shot PodcastIndex import
After=network-online.target stophammer-primary.service
Wants=network-online.target

[Service]
Type=oneshot
User=stophammer-crawler
Group=stophammer-crawler
EnvironmentFile=/etc/stophammer/crawler-import.env
ExecStart=/bin/sh -lc 'set -- /usr/bin/stophammer-crawler import --db "${PODCASTINDEX_DB:-/var/lib/stophammer-crawler/podcastindex_feeds.db}" --state "${IMPORT_STATE_DB:-/var/lib/stophammer-crawler/import_state.db}" --skip-db "${FEED_SKIP_DB:-/var/lib/stophammer-crawler/feed_skip.db}"; if [ -n "${CONCURRENCY:-}" ]; then set -- "$@" --concurrency "$CONCURRENCY"; fi; if [ -n "${IMPORT_BATCH:-}" ]; then set -- "$@" --batch "$IMPORT_BATCH"; fi; if [ -n "${IMPORT_CURSOR:-}" ]; then set -- "$@" --cursor "$IMPORT_CURSOR"; fi; if [ "${IMPORT_SKIP_KNOWN_NON_MUSIC:-0}" = "1" ]; then set -- "$@" --skip-known-non-music; fi; if [ "${IMPORT_SKIP_KNOWN_SUCCESS:-0}" = "1" ]; then set -- "$@" --skip-known-success; fi; exec "$@"'
WorkingDirectory=/var/lib/stophammer-crawler
StateDirectory=stophammer-crawler
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
ReadWritePaths=/var/lib/stophammer-crawler

[Install]
WantedBy=multi-user.target
